# vault.service
[Unit]
Description=Knowledge Vault API Server
After=network.target

[Service]
Type=simple
User=mrchartist
WorkingDirectory=/mnt/d/AG/Vault/apps/vault
ExecStart=/usr/bin/node /mnt/d/AG/Vault/apps/vault/server.mjs
Restart=always
RestartSec=5
StandardOutput=append:/mnt/d/AG/Vault/db/vault.log
StandardError=append:/mnt/d/AG/Vault/db/vault.log
Environment=VAULT_ROOT=/mnt/d/AG/Vault
Environment=VAULT_PORT=7799

[Install]
WantedBy=multi-user.target
